# deploy/safetyvision.service
[Unit]
Description=SafetyVision Real-Time Forklift Safety System
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=@RUN_USER@
Group=@RUN_GROUP@
SupplementaryGroups=audio video
WorkingDirectory=@TARGET_DIR@
ExecStart=@VENV_DIR@/bin/python -m safetyvision.supervisor -c @CONFIG_PATH@
Restart=always
RestartSec=2
CPUAffinity=0 1 2 3
LimitNOFILE=8192
StandardOutput=journal
StandardError=journal
SyslogIdentifier=safetyvision
Environment=SAFETYVISION_CONFIG=@CONFIG_PATH@
Environment=PYTHONPATH=@TARGET_DIR@/src

[Install]
WantedBy=multi-user.target
